# .config/systemd/user/autorotate.service
[Unit]
Description=autorotate

[Service]
ExecStart=/home/bou/app/bin/autorotate.sh

[Install]
WantedBy=default.target
